# skytest/main.cron
01 3 * * * /home/skyeng/tolmat/main.py
